# debian/tinyllama-mdns/etc/cron.d/tinyllama-mdns
# This fixes issues caused by avahi-daemon failing to process mDNS requests after 60s
*/1 * * * * root systemctl restart avahi-daemon.service
