# debian/tlweb-swarmui/DEBIAN/prerm
#!/bin/bash
systemctl stop swarmui ||:
systemctl disable swarmui ||: