# debian/tinyllama-openedai-speech/DEBIAN/prerm
#!/bin/bash
systemctl stop openedai-speech ||:
systemctl disable openedai-speech ||: