# debian/ollama/DEBIAN/prerm
#!/bin/bash
systemctl stop ollama ||:
systemctl disable ollama ||: