# debian/tlweb-stable-diffusion-webui/DEBIAN/prerm
#!/bin/bash
systemctl stop stable-diffusion-webui ||:
systemctl disable stable-diffusion-webui ||: