# debian/tlweb-open-webui/DEBIAN/prerm
#!/bin/bash
systemctl stop open-webui ||:
systemctl disable open-webui ||: