# debian/tinyllama/DEBIAN/prerm
#!/bin/bash
systemctl stop tinyllama ||:
systemctl disable tinyllama ||: