# debian/tlweb-amica/DEBIAN/prerm
#!/bin/bash
systemctl stop amica ||:
systemctl disable amica ||: