# debian/tlweb-stableswarmui/DEBIAN/prerm
#!/bin/bash
systemctl stop stableswarmui ||:
systemctl disable stableswarmui ||: