# debian/tlweb-langflow/DEBIAN/prerm
#!/bin/bash
systemctl stop langflow ||:
systemctl disable langflow ||: